# ansible/synology/templates/synology-observability.service
[Unit]
Description=Emit Synology backup health events to syslog
After=syslog-ng.service

[Service]
Type=oneshot
ExecStart=/usr/local/sbin/synology-observability-logger
